# roles/ansible-monasca-elkstack/templates/kibana.service.j2
[Unit]
Description=Kibana
Requires=network.target
After=network.target

[Service]
Type=simple
Restart=on-failure
ExecStart={{ kibana_dest }}/bin/kibana

[Install]
WantedBy=multi-user.target
